# config/systemd/user/backup-taskw.service
[Unit]
Description=Back up taskw db

[Service]
ExecStart=/usr/bin/bash %h/bin/backup-taskw.sh
