# config/systemd/user/offlineimap-quick.service
[Unit]
Description=Quick mail sync

[Service]
ExecStart=/usr/bin/bash /home/asinha/bin/fetch-mail.sh -q
