# config/systemd/user/offlineimap-full.service
[Unit]
Description=Full mail sync

[Service]
ExecStart=/usr/bin/bash %h/bin/fetch-mail.sh -f
