# config/systemd/user/task-report-daily.service
[Unit]
Description=Run task-report

[Service]
ExecStart=/usr/bin/bash %h/bin/make-task-report.sh -t
