# config/systemd/user/timew-notification.service
[Unit]
Description=Run time-notification

[Service]
ExecStart=/usr/bin/bash %h/bin/timew-notification.sh -n
